# deploy/ecab-octane.conf
; =============================================================
; ECAB — Laravel Octane (RoadRunner) Supervisor Configuration
; =============================================================
; File:   /etc/supervisor/conf.d/ecab-octane.conf
; Reload: sudo supervisorctl reread && sudo supervisorctl update
; Start:  sudo supervisorctl start ecab-octane:*
; Restart: sudo supervisorctl restart ecab-octane:*
; Logs:  /var/log/ecab/octane.log
; =============================================================

[program:ecab-octane]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/ecab/current/artisan octane:start
    --server=roadrunner
    --host=127.0.0.1
    --port=8000
    --workers=4
    --task-workers=2
    --max-requests=500
autostart=true
autorestart=true
; Give workers 30s to finish in-flight requests before killing
stopwaitsecs=30
stopasgroup=true
killasgroup=true
; Run as the web server user (matches Nginx/FPM)
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/ecab/octane.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/ecab/octane-error.log

; =============================================================
; ECAB — Queue Worker (for RetryPoolMatch + ProcessPoolTimeout)
; =============================================================
; Running multiple workers so pool retry jobs don't block each other
[program:ecab-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/ecab/current/artisan queue:work redis
    --sleep=3
    --tries=3
    --max-time=3600
    --queue=default,pool-retries
autostart=true
autorestart=true
stopwaitsecs=60
stopasgroup=true
killasgroup=true
user=www-data
; 4 queue workers = up to 4 concurrent pool retries
numprocs=4
redirect_stderr=true
stdout_logfile=/var/log/ecab/queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
